# supervisord.conf

[supervisord]
nodaemon=true
user=appuser

[program:uwsgi]
command=/opt/app/venv/bin/uwsgi --ini /opt/app/uwsgi.ini --die-on-term
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0